# systemctl/puma.staging.service
[Unit]
Description=Puma HTTP Server (staging)
After=network.target

[Service]
Type=simple
Environment=RAILS_ENV=staging
WorkingDirectory=/var/www/myapp/current
User=deploy
Group=deploy

ExecStart=/bin/bash -lc "bundle exec puma -C /var/www/myapp/shared/puma.rb"
ExecReload=/bin/kill -TSTP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

RestartSec=1
Restart=on-failure

StandardOutput=syslog
StandardError=syslog

SyslogIdentifier=puma

[Install]
WantedBy=default.target
